# .config/cron/crontab

# newsboat feed update
0 * * * * newsboat -x reload
